# tools/image-builder/files_debootstrap/etc/systemd/system/ssh-setup.service
[Unit]
Description=SSH setup for Firecracker VM
Before=ssh.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssh-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=firecracker.target